# packages/pglift/pglift-0.24.0.tar.gz/pglift-0.24.0/src/pglift/data/systemd/pglift-backup@.timer
[Unit]
Description=Backup %i PostgreSQL database instance
Requires=pglift-postgresql@%i.service
After=pglift-postgresql@%i.service

[Timer]
OnCalendar={calendar}
Persistent=true

[Install]
WantedBy=timers.target
