# packages/pglift/pglift-0.24.0-py3-none-any.whl/pglift/data/systemd/pglift-backup@.service
[Unit]
Description=Backup %i PostgreSQL database instance
After=pglift-postgresql@%i.service

[Service]
Type=oneshot
{executeas}

{environment}
ExecStart={python} -m pglift.backup %i
